# packaging/debian/media-dashboard.service
[Unit]
Description=Home Media Dashboard
After=network.target

[Service]
ExecStart=/usr/bin/media-dashboard
EnvironmentFile=-%h/.config/media-dashboard/env
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
